# tools/rfc2217_server/supervisor_rfc2217_server.conf
[program:rfc2217_server_tentsensors]
command=python /path/to/rfc2217_server.py -p 1234 /dev/ttyUSB0
directory=/path/to
environment=PATH="/path/to/python/virtualenv/bin"
autostart=true
autorestart=true
startretries=30
user=pi
